# rc.FreeBSD
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: launchd
# REQUIRE: DAEMON
# BEFORE:  LOGIN

. /etc/rc.subr

name="launchd"
rcvar="launchd_enable"
command="/usr/local/sbin/${name}"
start_postcmd="/usr/local/bin/launchctl load /usr/local/etc/launchd/daemons /usr/local/share/launchd/daemons"

load_rc_config $name
run_rc_command "$1"
