# rc.FreeBSD
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: jobd
# REQUIRE: DAEMON
# BEFORE:  LOGIN

. /etc/rc.subr

name="jobd"
rcvar="jobd_enable"
command="/usr/local/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
